# dist/phoebe.service
[Unit]
Description=phoebe automatic tuning
Documentation=man:phoebe(8)
After=network.target

[Service]
ExecStart=/usr/bin/phoebe -f /usr/share/phoebe/rates.csv -m inference -s  /etc/phoebe/settings.json
Type=simple

#Restart=on-failure

[Install]
WantedBy=multi-user.target
